# fuckrouter.service
[Unit]
Description=Fuck router

[Service]
Type=simple
ExecStart=fuckrouter
