# fuckrouter.service
[Unit]
Description=Fuck router

[Service]
ExecStart=fuckrouter

[Install]
WantedBy=network-online.target
